# trading/deploy/aws/systemd/sfo-kalshi-paper-settle.timer
[Unit]
Description=Run SFO Kalshi paper settlement checks every 30 minutes

[Timer]
# Decorrelate firing so the sfo units don't all spike together
RandomizedDelaySec=30
OnCalendar=*-*-* *:10,40
Persistent=true
Unit=sfo-kalshi-paper-settle.service

[Install]
WantedBy=timers.target
