# trading/deploy/aws/systemd/sfo-forecaster-refresh.service.in
[Unit]
Description=WeatherEdge multi-city forecaster refresh (Google/NWS/EMOS)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=__APP_USER__
WorkingDirectory=__FORECASTER_DIR__
EnvironmentFile=__ENV_FILE__
# The external forecast refresh is a legitimately long oneshot (~1-2 min);
# raise the start timeout above the 90 s default so a normal run is not killed.
TimeoutStartSec=900
ExecStart=/bin/bash __TRADING_DIR__/deploy/aws/sync_forecaster_source.sh
# External fetches get a '-' prefix: a transient NWS/Google hiccup must NOT abort
# the remaining refresh steps. The local re-blend runs on whatever cached data
# did land; freshness is still gated by max_forecast_age_hours and the watchdog.
# Station observations for every configured city (intraday context +
# observed-high conditioning). Two local days keeps it ~30 requests per tick.
ExecStart=-__FORECASTER_DIR__/.venv/bin/python nws_ground_truth.py --days 2 --cities all
# The Google/blend stack remains SFO-only (metered API budget).
ExecStart=-__FORECASTER_DIR__/.venv/bin/python google_weather_cache.py --refresh
ExecStart=__FORECASTER_DIR__/.venv/bin/python google_weather_cache.py
# Refresh station-keyed CLI truth before EMOS reads its training/serve history.
# Fail-soft: a transient product outage must not freeze live forecasts.
ExecStart=-__FORECASTER_DIR__/.venv/bin/python __FORECASTER_DIR__/city_truth.py --db __FORECASTER_DIR__/weather.db --refresh --cities all
# Live EMOS for every city's rolling window: ONE batched Open-Meteo request
# per city per tick. The NWP previous-runs archive and the rolling-origin EMOS
# rebuild moved to the nightly maintenance unit -- they are training data, not
# live state, and at fifteen cities a 30-minute cadence would blow the free
# API budget. '-' prefix => a fetch hiccup never fails the refresh pipeline.
ExecStart=-__FORECASTER_DIR__/.venv/bin/python emos_forecast.py --serve-rolling --cities all
Nice=5
# Contain a runaway to this unit instead of letting the kernel OOM-killer
# pick an arbitrary victim (4 GB box: containment, not thrash-protection).
# Tune after observing real peaks via `systemctl status`.
MemoryHigh=800M
MemoryMax=1100M
