# trading/deploy/aws/systemd/sfo-kalshi-paper-monitor.timer
[Unit]
Description=Run SFO Kalshi paper exit monitor every 2 minutes

[Timer]
# Decorrelate firing so the sfo units don't all spike together
RandomizedDelaySec=20
OnCalendar=*-*-* *:01,03,05,07,09,11,13,15,17,19,21,23,25,27,29,31,33,35,37,39,41,43,45,47,49,51,53,55,57,59
Persistent=true
Unit=sfo-kalshi-paper-monitor.service

[Install]
WantedBy=timers.target
