# trading/deploy/aws/systemd/sfo-dataset-backfill.timer
[Unit]
Description=Run SFO compact dataset backfill nightly

[Timer]
# Decorrelate firing so the sfo units don't all spike together
RandomizedDelaySec=120
OnCalendar=*-*-* 02:25:00
Persistent=true
Unit=sfo-dataset-backfill.service

[Install]
WantedBy=timers.target
