# trading/deploy/aws/systemd/sfo-kalshi-paper-scan.timer
[Unit]
Description=Run SFO Kalshi paper scans every 5 minutes

[Timer]
# Decorrelate firing so the sfo units don't all spike together
RandomizedDelaySec=30
# Every 5 minutes: the scan live-fetches the current Kalshi order books and
# makes paper-trade entries on fresh market data, so a newly-listed bracket (and
# the comfortable far-tail NO it may offer) is acted on within ~5 min instead of
# waiting up to 15. The 2-minute exit monitor and 5-minute dashboard refresh are
# unchanged.
OnCalendar=*-*-* *:00,05,10,15,20,25,30,35,40,45,50,55
Persistent=true
Unit=sfo-kalshi-paper-scan.service

[Install]
WantedBy=timers.target
