# trading/deploy/aws/systemd/sfo-forecast-freshness.timer
[Unit]
Description=Run SFO forecast freshness watchdog every 30 minutes

[Timer]
# Decorrelate firing; check shortly after each refresh window.
RandomizedDelaySec=60
OnCalendar=*-*-* *:15,45:00
Persistent=true
Unit=sfo-forecast-freshness.service

[Install]
WantedBy=timers.target
