# trading/deploy/aws/systemd/sfo-strategy-lab-refresh.timer
[Unit]
Description=Run SFO Strategy Lab trading-results refresh every fifteen minutes

[Timer]
# Decorrelate firing so the sfo units don't all spike together
RandomizedDelaySec=30
OnActiveSec=4min
OnUnitActiveSec=15min
AccuracySec=30s
Persistent=true
Unit=sfo-strategy-lab-refresh.service

[Install]
WantedBy=timers.target
