# trading/deploy/aws/systemd/sfo-kalshi-paper-monitor.service.in
[Unit]
Description=SFO Kalshi paper exit monitor
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=__APP_USER__
WorkingDirectory=__TRADING_DIR__
EnvironmentFile=__ENV_FILE__
ExecStart=__TRADING_DIR__/.venv/bin/python -m sfo_kalshi_quant.cli --no-color --forecaster-root __FORECASTER_DIR__ --db-path __TRADING_DIR__/data/paper_trading.db paper-monitor --take-profit-pct ${PAPER_TAKE_PROFIT_PCT} --stop-loss-pct ${PAPER_STOP_LOSS_PCT} --yes-take-profit-pct ${PAPER_YES_TAKE_PROFIT_PCT} --yes-stop-loss-pct ${PAPER_YES_STOP_LOSS_PCT} --no-take-profit-pct ${PAPER_NO_TAKE_PROFIT_PCT} --no-stop-loss-pct ${PAPER_NO_STOP_LOSS_PCT} --model-veto-max-loss-pct ${PAPER_MODEL_VETO_MAX_LOSS_PCT} --model-veto-buffer ${PAPER_MODEL_VETO_BUFFER}
Nice=5
# Contain a runaway process to this unit instead of letting the kernel
# OOM-killer pick an arbitrary victim (4 GB box: containment, not
# thrash-protection).
MemoryHigh=400M
MemoryMax=512M
