# trading/deploy/aws/systemd/sfo-operational-publish.timer
[Unit]
Description=Build and publish WeatherEdge operational artifacts every five minutes

[Timer]
# Stagger publication from the forecaster and paper-trading timers.
RandomizedDelaySec=30
OnActiveSec=2min
OnUnitActiveSec=5min
AccuracySec=30s
Persistent=true
Unit=sfo-operational-publish.service

[Install]
WantedBy=timers.target
