# trading/deploy/aws/systemd/sfo-kalshi-paper-scan.service.in
[Unit]
Description=SFO Kalshi paper trading scan
Wants=network-online.target
After=network-online.target
OnFailure=sfo-alert@%n.service

[Service]
Type=oneshot
User=__APP_USER__
WorkingDirectory=__TRADING_DIR__
EnvironmentFile=__ENV_FILE__
ExecStart=/usr/bin/env bash __TRADING_DIR__/deploy/aws/run_paper_scan_profiles.sh
Nice=5
# Contain a runaway process to this unit instead of letting the kernel
# OOM-killer pick an arbitrary victim (4 GB box: containment, not
# thrash-protection).
MemoryHigh=400M
MemoryMax=512M
