# trading/deploy/aws/systemd/sfo-dataset-backfill.service.in
[Unit]
Description=WeatherEdge nightly dataset + multi-city archive maintenance
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=__APP_USER__
WorkingDirectory=__TRADING_DIR__
EnvironmentFile=__ENV_FILE__
ExecStart=/usr/bin/env bash __TRADING_DIR__/deploy/aws/run_dataset_backfill.sh
# Retention for the ~60k-rows/day decision journal (see PaperStore.prune_decision_snapshots).
ExecStart=-__TRADING_DIR__/.venv/bin/python -m sfo_kalshi_quant.cli --no-color paper-prune
# Multi-city training-data maintenance (nightly is enough; the newest scored
# day lands within one day, immaterial against a ~340-day archive):
# yesterday's CLI settlements from the IEM archive (also picks up corrected
# finals), the leakage-free NWP previous-runs rows, and the rolling-origin
# EMOS rebuild that refreshes every city's calibration record.
ExecStart=-__FORECASTER_DIR__/.venv/bin/python __FORECASTER_DIR__/city_truth.py --db __FORECASTER_DIR__/weather.db --backfill-iem --start-year 2026 --cities all
ExecStart=-__FORECASTER_DIR__/.venv/bin/python __FORECASTER_DIR__/nwp_archive.py --db __FORECASTER_DIR__/weather.db --daily --cities all
ExecStart=-__FORECASTER_DIR__/.venv/bin/python __FORECASTER_DIR__/emos_forecast.py --db __FORECASTER_DIR__/weather.db --backfill --lead 1 --cities all
ExecStart=-__FORECASTER_DIR__/.venv/bin/python __FORECASTER_DIR__/emos_forecast.py --db __FORECASTER_DIR__/weather.db --backfill --lead 2 --cities all
Nice=10
# Per-unit memory cap so a spike self-limits on the no-swap 1 GB box.
MemoryHigh=280M
MemoryMax=360M
